# scripts/robot.service
[Unit]
Description=robot service

[Service]
User=andrew
WorkingDirectory=/home/andrew/snow/
ExecStart=/home/andrew/snow/snow
Restart=always

[Install]
WantedBy=multi-user.target
